# rootfs/etc/services.d/saslauthd/run
#!/usr/bin/with-contenv bash
exec saslauthd -a "${SASL_AUTH_MECHANISM:-ldap}" -r -d -O /etc/saslauthd/saslauthd.conf
